# packaging/systemd/polkit-agent-helper@.service.d/sentinel.conf
# Drop-in for polkit-agent-helper@.service shipped by Sentinel.
#
# The vendor unit sets ProtectHome=yes which masks /run/user/<uid>
# inside the sandbox. Sentinel's bypass socket lives at
# $XDG_RUNTIME_DIR/sentinel-agent.sock (= /run/user/<uid>/sentinel-agent.sock),
# so without this override pam_sentinel.so inside helper-1 can't reach
# the agent and every polkit-mediated auth would re-prompt for a
# password.
#
# Disabling ProtectHome is harmless for helper-1: it still runs under
# all the other sandbox restrictions (PrivateTmp, ProtectSystem=strict,
# RestrictAddressFamilies=AF_UNIX, MemoryDenyWriteExecute, NoNewPrivileges,
# SystemCallFilter=@system-service). Helper-1's job is to talk to PAM
# and the polkitd D-Bus channel — neither needs /home or /root.

[Service]
ProtectHome=no
