# packaging/xdg-autostart/sentinel-polkit-agent.desktop
# SPDX-FileCopyrightText: 2025 Atay Özcan <atay@oezcan.me>
# SPDX-License-Identifier: GPL-3.0-or-later
[Desktop Entry]
Type=Application
Name=Sentinel Polkit Authentication Agent
Comment=Authenticates you to perform privileged tasks via the Sentinel UAC-style dialog
Exec=/usr/lib/sentinel-polkit-agent
Categories=System;
NoDisplay=true
X-GNOME-Autostart-Phase=Initialization
# Critical: tell systemd-xdg-autostart-generator NOT to wrap this entry.
# When systemd wraps it as `app-*@autostart.service`, the agent runs as a
# child of `user@.service` and inherits the manager-session sessionid — so
# polkit's RegisterAuthenticationAgent rejects it with "Passed session and
# the session the caller is in differs". By skipping the systemd wrapper
# we let cosmic-session / Hyprland / Sway autostart spawn the agent as a
# direct child of the compositor, inheriting the graphical sessionid.
#
# `X-GNOME-Autostart-enabled=false` and `X-systemd-skip=true` together
# cover both old and new versions of the systemd generator — the former
# is the GNOME-style key the generator has long honoured; the latter is
# the documented systemd-native mechanism.
X-GNOME-Autostart-enabled=false
X-systemd-skip=true
# Polkit only allows ONE agent per session — exclude desktops that
# already ship their own. Everything else (cosmic-comp, Hyprland, sway,
# Niri, River, Wayfire, …) gets the agent autostarted; new wlroots-style
# compositors are picked up automatically without needing to be added
# to a positive allowlist.
NotShowIn=GNOME;KDE;XFCE;LXDE;Cinnamon;MATE;LXQt;Pantheon;Budgie;